# rootfs/etc/services.d/zigbee2mqtt/run
#!/usr/bin/with-contenv bash

#####################################################

# CMD=$1
cd /app



exec npm start 
